# recipes-support/bootstrap-autorun/files/bootstrap-autorun.path
[Path]
PathExists=/dev/disk/by-label/BOOTSTRAP

[Install]
WantedBy=multi-user.target
