# recipes-core/base-files/base-files/data.mount
[Unit]
Description=Mount DATA partition on /data

[Mount]
What=/dev/disk/by-label/DATA
Where=/data
