# recipes-core/base-files/base-files/kubelet.path
[Unit]
Description="Watch for kubeconfig"

[Path]
PathExists=/etc/kubernetes/kubelet.conf

[Install]
WantedBy=multi-user.target
