# recipes-core/base-files/base-files/var-lib.mount
[Unit]
Description=DATA partition on /var/lib

[Mount]
What=/dev/disk/by-label/DATA
Where=/var/lib
